# packages/ipnoise-router/rootfs/lib/systemd/system/ipnoise-router.service
[Unit]
Description=IPNoise router
After=syslog.target network.target

[Service]
ExecStart=/usr/sbin/ipnoise-router
Restart=always
IgnoreSIGPIPE=no

[Install]
WantedBy=multi-user.target

