# packages/ipnoise-hostfs/rootfs/lib/systemd/system/ipnoise-hostos.service
[Unit]
Description=IPNoise hostfs support
After=syslog.target network.target

[Service]
ExecStart=/usr/sbin/ipnoise-hostfs
Restart=always
IgnoreSIGPIPE=no

[Install]
WantedBy=multi-user.target

