# packages/ipnoise-squid/root/lib/systemd/system/ipnoise-squid.service
[Unit]
Description=IPNoise squid proxy server
After=syslog.target network.target

[Service]
ExecStart=/usr/sbin/squid -N
Restart=always
IgnoreSIGPIPE=no

[Install]
WantedBy=multi-user.target

