# packages/ipnoise-transport/rootfs/lib/systemd/system/ipnoise-transport-udp-v4.service
[Unit]
Description=IPNoise UDP v4. transport
After=syslog.target network.target

[Service]
ExecStart=/usr/sbin/ipnoise-transport-udp-v4
Restart=always
IgnoreSIGPIPE=no

[Install]
WantedBy=multi-user.target

